# wrappers/systemd.template.service.in
# do not edit this file in /lib/systemd/system - instead, do the following:
# cp /lib/systemd/system/dirsrv\@.service /etc/systemd/system/dirsrv\@.service
# mkdir -p /etc/systemd/system/@systemdgroupname@.wants
# edit /etc/systemd/system/dirsrv\@.service - uncomment the LimitNOFILE=8192 line
# where %i is the name of the instance
# you may already have a symlink in
# /etc/systemd/system/@systemdgroupname@.wants/dirsrv@%i.service pointing to
# /lib/systemd/system/dirsrv\@.service - you will have to change it to link
# to /etc/systemd/system/dirsrv\@.service instead
# ln -s /etc/systemd/system/dirsrv\@.service /etc/systemd/system/@systemdgroupname@.wants/dirsrv@%i.service
# systemctl daemon-reload 
# systemctl (re)start @systemdgroupname@
[Unit]
Description=@capbrand@ Directory Server %i.
BindTo=@systemdgroupname@
After=@systemdgroupname@

[Service]
Type=forking
EnvironmentFile=@initconfigdir@/@package_name@
EnvironmentFile=@initconfigdir@/@package_name@-%i
ExecStart=@sbindir@/ns-slapd -D @instconfigdir@/slapd-%i -i @localstatedir@/run/@package_name@/slapd-%i.pid -w @localstatedir@/run/@package_name@/slapd-%i.startpid
ExecStopPost=/bin/rm -f @localstatedir@/run/@package_name@/slapd-%i.pid
# uncomment this line to raise the file descriptor limit - but do not edit this
# file in /lib/systemd/system - see above for instructions
# LimitNOFILE=8192
